# examples/powerplant/source/cpp/powerplant
#
# Copyright 2014-2015 CyberVision, Inc.
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#      http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#

#! /bin/bash

### BEGIN INIT INFO
# Provides:          powerplant
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

# NOTE: The $APP_NAME value should match the one in the build.sh script.
APP_NAME="powerplant_demo"

case "$1" in
    start)
        echo "Starting $APP_NAME service..."
        $APP_NAME &
        ;;
    stop)
        echo "Stopping $APP_NAME service..."
        killall $APP_NAME
        ;;
    restart)
        echo "Restarting $APP_NAME service..."
        killall $APP_NAME
        $APP_NAME &
        ;;
esac